# copy/tmp/mdata-setup.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type='manifest' name='mdata-setup'>
<service name='coreio/mdata-setup' type='service' version='1'>
        <create_default_instance enabled='false' />
        <single_instance />
        <dependency name="network" grouping="require_all" restart_on="error" type="service">
                <service_fmri value="svc:/milestone/network:default"/>
        </dependency>
        <dependency name='fs-root' grouping='require_all' restart_on='none' type='service'>
                <service_fmri value='svc:/system/filesystem/root' />
        </dependency>
        <exec_method type='method' name='start' exec='/opt/core/bin/mdata-setup' timeout_seconds='0' />
        <exec_method type='method' name='stop'  exec=':true' timeout_seconds='0' />
        <property_group name='startd' type='framework'>
                <propval name='duration' type='astring' value='transient' />
        </property_group>
        <stability value='Unstable' />
        <template>
            <common_name>
                <loctext xml:lang="C">Apply mdata values to configuration</loctext>
            </common_name>
        </template>
</service>
</service_bundle>
