# copy/opt/core/lib/svc/manifest/spamd.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='network/spamd' type='service' version='0'>
    <create_default_instance enabled='false'/>
    <single_instance/>
    <dependency name='network' grouping='require_all' restart_on='error' type='service'>
      <service_fmri value='svc:/milestone/network:default'/>
    </dependency>
    <dependency name='filesystem-local' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local:default'/>
    </dependency>
    <dependency name='autofs' grouping='optional_all' restart_on='error' type='service'>
      <service_fmri value='svc:/system/filesystem/autofs:default'/>
    </dependency>
    <exec_method name='start' type='method' exec='/opt/local/bin/spamd -r /var/spamassassin/spamd.pid -i 127.0.0.1 -p 11333 -u spamd --syslog-socket=none -s /var/log/spamd/spamd.log' timeout_seconds='60'>
    <method_context working_directory='/var/spamassassin' project='default'>
      <method_credential user='root' group='root'/>
      <method_environment>
        <envvar name='PATH' value='/opt/local/bin:/opt/local/sbin:/usr/bin:/usr/sbin'/>
        <envvar name='TMPDIR' value='/var/tmp'/>
      </method_environment>
    </method_context>
    </exec_method>
    <exec_method name='stop' type='method' exec=':kill' timeout_seconds='60'>
    <method_context working_directory='/var/spamassassin' project='default'>
      <method_credential user='spamd' group='spamd'/>
      <method_environment>
        <envvar name='HOME' value='/var/spamassassin'/>
        <envvar name='PATH' value='/opt/local/bin:/opt/local/sbin:/usr/bin:/usr/sbin'/>
        <envvar name='TMPDIR' value='/var/tmp'/>
      </method_environment>
    </method_context>
    </exec_method>
    <property_group name='startd' type='framework'>
      <propval name='duration' type='astring' value='child'/>
      <propval name='ignore_error' type='astring' value='core,signal'/>
    </property_group>
    <stability value='Evolving'/>
    <template>
      <common_name>
        <loctext xml:lang='C'>Spam Assassin</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
